# scripts/S99main
#!/bin/sh /etc/rc.common

START=99
STOP=10

PROG="/usr/main"
PIDFILE="/var/run/main.pid"
LOGFILE="/tmp/main.log"

start() {
    if [ ! -x "$PROG" ]; then
        echo "[S99main] missing $PROG" >&2
        return 1
    fi

    stop

    : > "$LOGFILE"
    "$PROG" >>"$LOGFILE" 2>&1 &
    pid=$!
    echo "$pid" > "$PIDFILE"

    sleep 1
    if ! kill -0 "$pid" 2>/dev/null; then
        echo "[S99main] $PROG exited early, see $LOGFILE" >&2
        rm -f "$PIDFILE"
        return 1
    fi
}

stop() {
    if [ -f "$PIDFILE" ]; then
        pid=$(cat "$PIDFILE" 2>/dev/null)
        if [ -n "$pid" ]; then
            kill "$pid" 2>/dev/null || true
        fi
        rm -f "$PIDFILE"
    fi

    killall main 2>/dev/null || true
}
